# snes-adb.service
[Unit]
Description=SNES ADB Daemon (Android TV remote)
After=network-online.target snes-controller.service
Wants=network-online.target

[Service]
EnvironmentFile=/home/wil/.env
ExecStart=/usr/bin/python3 /home/wil/snes_adb_daemon.py
Restart=always
RestartSec=5
User=wil
Group=wil
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
